# usr/lib/systemd/system/chupa-text.service
[Unit]
Description=ChupaText - Text extraction service
After=network.target docker.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/lib/chupa-text
ExecStart=/usr/local/bin/docker-compose up --abort-on-container-exit
ExecStop=/usr/local/bin/docker-compose down
# Because docker-compose up --abort-on-container-exit returns 0. :<
Restart=always
TimeoutSec=5min

[Install]
WantedBy=multi-user.target
